# config/systemd/user/mpdas.service
[Unit]
Description=MPD Scribbler
After=mpd.service

[Service]
ExecStart=/home/rev/src/mpdas/mpdas

[Install]
WantedBy=default.target
