# config/systemd/user/mpdas.service
[Unit]
Description=MPD Scribbler
After=mpd.service
Requires=mpd.service

[Service]
ExecStart=/usr/bin/mpdas

[Install]
WantedBy=default.target
