# config/systemd/user/yams.service.d/override.conf
[Unit]
Requires=mpd.service

[Service]
ExecStart=
ExecStart=/usr/bin/python3 -m yams -N -c %h/.cache/yams.cache
Restart=on-failure
RestartSec=5
